# mydeb/etc/systemd/system/hd.service
[Unit]
Description=HD Service
After=network.target

[Service]
ExecStart=/usr/local/bin/hd
Restart=always

[Install]
WantedBy=multi-user.target

