# mydeb/etc/systemd/system/oss.service
[Unit]
Description=OSS Service
After=network.target

[Service]
ExecStart=/usr/local/bin/oss
Restart=always

[Install]
WantedBy=multi-user.target
